# orchestratia-agent.service
[Unit]
Description=Orchestratia Agent Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/orchestratia-agent
ExecStart=/usr/bin/python3 /opt/orchestratia-agent/daemon.py --config /etc/orchestratia/config.yaml
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orchestratia-agent

# Note: No filesystem restrictions (ProtectSystem, ReadWritePaths, etc.)
# because PTY child processes (Claude Code) need full filesystem access
# to the working directories they operate in.

[Install]
WantedBy=multi-user.target
